# files/etc/sv/grow-ssd/run
#!/bin/sh
# grow ssd partition and filesystem at boot (one-shot)

set -eu

# Wait for /ssd to be mounted (up to 60 seconds)
MAX_WAIT=60
WAITED=0
SSD_DEV=$(findmnt -n -o SOURCE /ssd 2>/dev/null || echo "")

case "$SSD_DEV" in
    /dev/*)
        # Already mounted, proceed
        ;;
    *)
        # Not yet mounted, wait for it
        echo "grow-ssd: waiting for /ssd to be mounted..."
        while [ $WAITED -lt $MAX_WAIT ]; do
            SSD_DEV=$(findmnt -n -o SOURCE /ssd 2>/dev/null || echo "")
            case "$SSD_DEV" in
                /dev/*)
                    echo "grow-ssd: /ssd is mounted on $SSD_DEV"
                    break
                    ;;
            esac
            sleep 1
            WAITED=$((WAITED + 1))
        done
        # Check if we actually found the device
        case "$SSD_DEV" in
            /dev/*)
                # Successfully detected, continue
                ;;
            *)
                echo "grow-ssd: unable to detect source device for /ssd; mount point /ssd may not exist or may not be ready yet; skipping"
                sv down /var/service/grow-ssd >/dev/null 2>&1 || true
                exit 0
                ;;
        esac
        ;;
esac

PKNAME=$(lsblk -no pkname "$SSD_DEV" 2>/dev/null || echo "")
if [ -n "$PKNAME" ]; then
    DISK="/dev/$PKNAME"
else
    case "$SSD_DEV" in
        # Bare NVMe namespace device, e.g. /dev/nvme0n1
        /dev/nvme[0-9]n[0-9])
            DISK="$SSD_DEV"
            ;;
        # NVMe partition device, e.g. /dev/nvme0n1p1 -> /dev/nvme0n1
        /dev/nvme[0-9]n[0-9]p[0-9]*)
            DISK=$(printf '%s\n' "$SSD_DEV" | sed -E 's/p[0-9]+$//')
            ;;
        # Non-NVMe devices with numeric partition suffix, e.g. /dev/sda1 -> /dev/sda
        /dev/*[0-9])
            DISK=$(printf '%s\n' "$SSD_DEV" | sed -E 's/[0-9]+$//')
            ;;
        # Fallback: leave as-is if we cannot reliably parse
        *)
            DISK="$SSD_DEV"
            ;;
    esac
fi

# Extract partition number from the device path
case "$SSD_DEV" in
    /dev/nvme*p[0-9]*)
        PART=$(printf '%s\n' "$SSD_DEV" | sed -E 's/.*p([0-9]+)$/\1/')
        ;;
    /dev/*[0-9])
        PART=$(printf '%s\n' "$SSD_DEV" | sed -E 's/.*[^0-9]([0-9]+)$/\1/')
        ;;
    *)
        PART=""
        ;;
esac

# No support for bare non-partitioned devices currently
case "$PART" in
    ''|*[!0-9]*)
        echo "grow-ssd: unable to determine partition number from $SSD_DEV or no partition; skipping"
        sv down /var/service/grow-ssd >/dev/null 2>&1 || true
        exit 0
        ;;
esac

# Validate device paths before proceeding with potentially destructive operations
# Ensure DISK matches expected block device patterns
case "$DISK" in
    /dev/sd[a-z]|/dev/sd[a-z][a-z]|/dev/nvme[0-9]*n[0-9]*|/dev/mmcblk[0-9]*|/dev/vd[a-z])
        # Valid block device pattern
        ;;
    *)
        echo "grow-ssd: device path '$DISK' does not match expected block device patterns; skipping for safety"
        sv down /var/service/grow-ssd >/dev/null 2>&1 || true
        exit 0
        ;;
esac

# Ensure DISK exists and is a block device
if ! [ -b "$DISK" ]; then
    echo "grow-ssd: '$DISK' is not a block device; skipping"
    sv down /var/service/grow-ssd >/dev/null 2>&1 || true
    exit 0
fi

# Ensure SSD_DEV exists and is a block device
if ! [ -b "$SSD_DEV" ]; then
    echo "grow-ssd: '$SSD_DEV' is not a block device; skipping"
    sv down /var/service/grow-ssd >/dev/null 2>&1 || true
    exit 0
fi

# Validate that PART is a reasonable partition number (1-128)
# Note: PART is already validated to be numeric at line 82, but we add range check for safety
if [ "$PART" -lt 1 ] 2>/dev/null || [ "$PART" -gt 128 ] 2>/dev/null; then
    echo "grow-ssd: partition number $PART is out of valid range (1-128); skipping"
    sv down /var/service/grow-ssd >/dev/null 2>&1 || true
    exit 0
fi

# Only run if partition can be grown
growpart --dry-run "$DISK" "$PART" >/dev/null 2>&1 || {
    sv down /var/service/grow-ssd >/dev/null 2>&1 || true
    exit 0
}

echo "Growing partition $SSD_DEV"
growpart "$DISK" "$PART"

echo "Resizing filesystem on $SSD_DEV"
resize2fs "$SSD_DEV"

echo "grow-ssd completed"

# Disable the service after successful completion
sv down /var/service/grow-ssd >/dev/null 2>&1 || true
